# systemd/whatsapp-engine.service
[Unit]
Description=Whatsapp Engine
After=network.target openclaw-whatsapp.service
Requires=openclaw-whatsapp.service

[Service]
Type=simple
WorkingDirectory=/root/.openclaw/workspace/Whatsapp Engine
ExecStart=/usr/bin/node /root/.openclaw/workspace/Whatsapp\x20Engine/src/server.js
Restart=always
RestartSec=5

[Install]
WantedBy=default.target
